# recipes-core/boot-mark-good/files/boot-mark-good.timer
[Unit]
Description=Timer to mark boot as good after 30 seconds

[Timer]
# Wait 30 seconds after boot before marking as good
OnBootSec=30
Unit=boot-mark-good.service

[Install]
WantedBy=timers.target
